# deploy/systemd/rwoo-closing-quotes-near.timer
[Unit]
Description=Capture targeted final-hour market quotes every five minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
RandomizedDelaySec=10
Persistent=true

[Install]
WantedBy=timers.target

